# Quadlet/aria2.container
[Unit]
Description=aria2-ui Quadlet

[Container]
Image=docker.io/hurlenko/aria2-ariang
ContainerName=aria2-ui
AutoUpdate=registry
Environment=TZ=Asia/Shanghai
Environment=ARIA2RPCPORT=443
Environment=RPC_SECRET=*****secret*****
Volume=/opt/data/download:/aria2/data
Volume=/opt/podman/aria2:/aria2/conf
Network=pasta:-i,br0
PublishPort=127.0.10.3:9999:8080

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
