# systemd/DDNS.service
[Unit]
Description=Execute DDNS script

[Service]
Type=oneshot
ExecStart=/opt/script/cf_ddns.sh ***cf_token***

[Install]
WantedBy=default.target
