# systemd/DDNS.timer
[Unit]
Description=Run DDNS script every 5 minutes

[Timer]
OnUnitActiveSec=5min
Unit=DDNS.service

[Install]
WantedBy=timers.target
