# systemd/DDNS.service
[Unit]
Description=Execute DDNS script

[Service]
Type=oneshot
ExecStart=/opt/script/DDNS.sh

[Install]
WantedBy=default.target
