# systemd/update-cert.service
[Unit]
Description=Update SSL Certificates

[Service]
Type=oneshot
ExecStart=%h/.local/bin/update-cert.sh
