# systemd/vaultwarden-backup.service
[Unit]
Description=Run Vaultwarden Backup Script

[Service]
Type=oneshot
ExecStart=%h/backups/vaultwarden-backup.sh --backup
