# etc/systemd/system/setsebool.service
[Service]
Type=oneshot 
ExecStart=setsebool container_use_dri_devices 1
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
